# themes/luci-theme-hashland/root/etc/uci-defaults/30_luci-theme-hashland
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Hashland=/luci-static/hashland
	set luci.main.mediaurlbase=/luci-static/hashland
	commit luci
EOF
exit 0
